# systemd/hostsblock.timer
[Unit]
Description=Daily hosts blocking update

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=network-online.target
